# bspwm/.config/systemd/user/lemonbar.socket
[Unit]
Description=lemonbar socket
Documentation=man:lemonbar(1)

[Socket]
ListenFIFO=%t/lemonbar-fifo

[Install]
WantedBy=sockets.target
