# bspwm/.config/systemd/user/lightstatus.socket
[Unit]
Description=lightstatus socket

[Socket]
ListenFIFO=%t/lightstatus-fifo

[Install]
WantedBy=sockets.target
