# jetbrains/.config/systemd/user/intellij-idea-community-edition.service
[Unit]
Description=IntelliJ IDEA Community Edition
PartOf=graphical-session.target
After=graphical-session-pre.target

[Service]
ExecStart=/usr/bin/idea.sh

[Install]
WantedBy=graphical-session.target
