# prime/.config/systemd/user/steam.service.d/prime-run.conf
[Service]
# Clear out the list first as otherwise we would add to it.
ExecStart=
# Run with prime-run so that any games are launched using dedicated graphics.
ExecStart=/usr/bin/prime-run steam
