# systemd/.config/systemd/user/redshift.service.d/override.conf
[Unit]
Requisite=display-server.target
After=display-server.target

[Install]
WantedBy=desktop.target
