# systemd/.config/systemd/user/bspwm.service
[Unit]
Description=bspwm - Binary space partitioning window manager
Documentation=man:bspwm(1)
Requires=sxhkd.service
Wants=bspwm-panel.service
After=display-server.target

[Service]
ExecStart=/usr/bin/bspwm

[Install]
Alias=window-manager.service
