# systemd/.config/systemd/user/bspwm.service
[Unit]
Description=bspwm - Binary space partitioning window manager
Requires=sxhkd.service
After=xorg.target

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/bspwm

[Install]
Also=bspwm-panel.service
WantedBy=wm.target
