# systemd/.config/systemd/user/bspwm.service
[Unit]
Description=bspwm
Documentation=man:bspwm(1)

Requires=sxhkd.service
Wants=bspwm-panel.service

Requires=x11.target
After=x11.target

[Service]
ExecStart=/usr/bin/bspwm

[Install]
Alias=window-manager.service
Also=x11.target
