# systemd/.config/systemd/user/desktop.target
[Unit]
Description=Desktop Target
Wants=window-manager.service

[Install]
Alias=default.target
