# systemd/.config/systemd/user/desktop.target
[Unit]
Description=Desktop (target)
Requires=xorg.target
Requires=window-manager.target

[Install]
Alias=default.target
