# systemd/.config/systemd/user/display-server.target
[Unit]
Description=Display Server (target)
