# systemd/.config/systemd/user/dunst.service
[Unit]
Description=dunst - A customizable and lightweight notification-daemon
After=dbus.socket
After=xorg.target

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=xinit.target
