# systemd/.config/systemd/user/sxhkd.service
[Unit]
Description=sxhkd - Simple X hotkey daemon
After=xorg.target

[Service]
ExecStart=/usr/bin/sxhkd
