# systemd/.config/systemd/user/sxhkd.service
[Unit]
Description=sxhkd
Documentation=man:sxhkd(1)
Requisite=display-server.target
After=display-server.target

[Service]
ExecStart=/usr/bin/sxhkd
