# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Udiskie

[Service]
ExecStart=/usr/bin/udiskie --verbose

[Install]
WantedBy=desktop.target
