# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie
Documentation=man:udiskie(8)

[Service]
ExecStart=/usr/bin/udiskie --verbose

[Install]
WantedBy=desktop.target
