# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Udiskie

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
