# systemd/.config/systemd/user/unclutter.service
[Unit]
Description=unclutter
Documentation=man:unclutter(1)
Requires=x11.target
After=x11.target

[Service]
ExecStart=/usr/bin/unclutter -noevents

[Install]
WantedBy=desktop.target
