# systemd/.config/systemd/user/unclutter.service
# Note - this is X-specific; if we were using wayland then this would fail.

[Unit]
Description=unclutter - remove idle cursor image from screen
After=display-server.target

[Service]
Type=simple
ExecStart=/usr/bin/unclutter -noevents

[Install]
WantedBy=desktop.target
