# systemd/.config/systemd/user/urxvtd.service
[Unit]
Description=RXVT-Unicode Daemon
After=display-server.target

[Service]
ExecStart=/usr/bin/urxvtd --opendisplay
Restart=on-failure
