# systemd/.config/systemd/user/urxvtd.service
[Unit]
Description=RXVT-Unicode Daemon
After=xorg.target

[Service]
Type=simple
ExecStart=/usr/bin/urxvtd -q -o

[Install]
WantedBy=xinit.target
