# systemd/.config/systemd/user/urxvtd.service
[Unit]
Description=RXVT-Unicode Daemon
Requires=urxvtd.socket
Requires=display-server.target
After=display-server.target

[Service]
ExecStart=/usr/bin/urxvtd --opendisplay
Restart=on-failure

[Install]
Also=urxvtd.socket
WantedBy=desktop.target
