# systemd/.config/systemd/user/window-manager.target
[Unit]
Description=Window Manager (target)
Requires=xorg.target
Requires=dbus.socket
# AllowIsolate=true
