# systemd/.config/systemd/user/window-manager.target
[Unit]
Description=Window Manager (target)
Requires=display-server.target
Requires=dbus.socket
