# systemd/.config/systemd/user/x11.socket
[Unit]
Description=X11 Server Socket

[Socket]
ListenStream=/tmp/.X11-unix/X0
Service=xorg.service

[Install]
WantedBy=sockets.target
