# systemd/.config/systemd/user/xorg.socket
[Unit]
Description=Xorg server (Socket)
Before=display-server.target

[Socket]
ListenStream=/tmp/.X11-unix/X0

[Install]
WantedBy=sockets.target
