# systemd/.config/systemd/user/xrdb.service
[Unit]
Description=xrdb - X server resource database utility
Documentation=man:xrdb(1)
After=x11.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/xrdb -load %h/.Xresources

[Install]
# This is X-specific; we want it only if X is used.
WantedBy=x11.target
