# systemd/.config/systemd/user/xrdb.service
[Unit]
Description=xrdb
Documentation=man:xrdb(1)
Requires=x11.target
After=x11.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/xrdb -merge %h/.Xresources

[Install]
WantedBy=desktop.target
