# systemd/.config/systemd/user/xrdb.service
[Unit]
Description=xrdb - X server resource database utility
Documentation=man:xrdb(1)
After=xorg.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/xrdb -load %h/.Xresources

[Install]
WantedBy=desktop.target
