# docker/.config/systemd/user/docker-image-prune@.timer
[Unit]
Description=Prune docker images weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
